# software/openwrt/attitude_adjustement/target/linux/ar71xx/linkduino/base-files/etc/uci-defaults/gpio-button-hotplug
#!/bin/sh
#
# Copyright (C) 2012 Michel Stempin
#

sed -i 's/\^button/button/' /etc/hotplug2.rules
